# playbooks/roles/backup/files/attic.cron
#!/bin/bash

/usr/local/bin/attic-backup backup
